# alice/iot/time_machine/misc/supervisord_deploy.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=2

[program:tvm_daemon]
command = /usr/local/bin/tvm.sh
priority=1
autorestart = true
stderr_logfile=/var/log/supervisor/tvm_daemon.err
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/supervisor/tvm_daemon.out
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2

[program:time_machine_server]
command=/usr/local/bin/timemachine
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:push_client]
command = /usr/local/bin/push-client.sh
priority=10
autorestart = true
redirect_stderr = true
stdout_logfile=/var/log/supervisor/push-client.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
