# jupytercloud/infrastructure/images/deploy-app/includes/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
user=root
autostart=true
startsecs=0
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO
directory=/
